# deployment/docker/services/db-backup/crontab
* * * * * PATH=/usr/local/bin:$PATH /entrypoint.sh >> /var/log/cron.log 2>&1
